# deploy/systemd/lobster-price-monitor-serve.service
[Unit]
Description=Lobster Price Monitor board HTTP server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=erik
WorkingDirectory=/opt/lobster-price-monitor
ExecStart=/opt/lobster-price-monitor/.venv/bin/python scripts/serve_board.py --port 8765 --host 0.0.0.0
Restart=on-failure
RestartSec=5
StandardOutput=append:/opt/lobster-price-monitor/logs/serve.log
StandardError=append:/opt/lobster-price-monitor/logs/serve.err

[Install]
WantedBy=multi-user.target
