# deploy/systemd/lobster-price-monitor-scrape.service
[Unit]
Description=Lobster Price Monitor scrape (no Telegram alerts)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=erik
WorkingDirectory=LOBSTER_ROOT
ExecStart=LOBSTER_ROOT/.venv/bin/python scripts/scrape_markets.py --no-alerts
StandardOutput=append:LOBSTER_ROOT/logs/scrape.stdout.log
StandardError=append:LOBSTER_ROOT/logs/scrape.stderr.log

[Install]
WantedBy=multi-user.target
