# deploy/systemd/lobster-price-monitor-scrape.service
[Unit]
Description=Lobster Price Monitor scrape (no Telegram alerts)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=erik
WorkingDirectory=/opt/lobster-price-monitor
ExecStart=/opt/lobster-price-monitor/.venv/bin/python scripts/scrape_markets.py --no-alerts
StandardOutput=append:/opt/lobster-price-monitor/logs/scrape.stdout.log
StandardError=append:/opt/lobster-price-monitor/logs/scrape.stderr.log

[Install]
WantedBy=multi-user.target
