# deploy/launchd/com.erik.lobster-price-monitor.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.erik.lobster-price-monitor.serve</string>

  <key>ProgramArguments</key>
  <array>
    <string>LOBSTER_ROOT/.venv/bin/python</string>
    <string>LOBSTER_ROOT/scripts/serve_board.py</string>
    <string>--port</string>
    <string>8765</string>
    <string>--host</string>
    <string>0.0.0.0</string>
  </array>

  <key>WorkingDirectory</key>
  <string>LOBSTER_ROOT</string>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>LOBSTER_ROOT/logs/serve.log</string>
  <key>StandardErrorPath</key>
  <string>LOBSTER_ROOT/logs/serve.err</string>
</dict>
</plist>
