# deploy/launchd/com.erik.lobster-price-monitor.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.erik.lobster-price-monitor.health</string>

  <key>ProgramArguments</key>
  <array>
    <string>LOBSTER_ROOT/.venv/bin/python</string>
    <string>LOBSTER_ROOT/scripts/health_check.py</string>
    <string>--log</string>
  </array>

  <key>WorkingDirectory</key>
  <string>LOBSTER_ROOT</string>

  <!-- Daily health snapshot at 06:00 local -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>6</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>LOBSTER_ROOT/logs/health.log</string>
  <key>StandardErrorPath</key>
  <string>LOBSTER_ROOT/logs/health.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
