# deploy/lobster-price-monitor-scrape.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.erik.lobster-price-monitor.scrape</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/erik-darci/Documents/Codex/2026-07-04/pu/work/lobster-price-monitor/.venv/bin/python</string>
    <string>/Users/erik-darci/Documents/Codex/2026-07-04/pu/work/lobster-price-monitor/scripts/scrape_markets.py</string>
    <string>--no-alerts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/erik-darci/Documents/Codex/2026-07-04/pu/work/lobster-price-monitor</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/erik-darci/Documents/Codex/2026-07-04/pu/work/lobster-price-monitor/logs/scrape.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/erik-darci/Documents/Codex/2026-07-04/pu/work/lobster-price-monitor/logs/scrape.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
