# deploy/lobster-price-monitor-serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.erik.lobster-price-monitor.serve</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/erik-darci/Documents/Codex/2026-07-04/pu/work/lobster-price-monitor/.venv/bin/python</string>
    <string>/Users/erik-darci/Documents/Codex/2026-07-04/pu/work/lobster-price-monitor/scripts/serve_board.py</string>
    <string>--port</string>
    <string>8765</string>
    <string>--bind</string>
    <string>0.0.0.0</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/erik-darci/Documents/Codex/2026-07-04/pu/work/lobster-price-monitor</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/erik-darci/Documents/Codex/2026-07-04/pu/work/lobster-price-monitor/logs/serve.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/erik-darci/Documents/Codex/2026-07-04/pu/work/lobster-price-monitor/logs/serve.stderr.log</string>
</dict>
</plist>
